# templates/default/dataloop-agent.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:          dataloop-agent
# Required-Start:    $local_fs $network
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: dataloop-agent
# Description:       dataloop monitoring agent
### END INIT INFO

start () {
  echo "dataloop agent started"
  /usr/local/bin/dataloop-lin-agent \
      --api-key <%= @api_key %> \
      --tags    <%= @tags %> \
      > /dev/null 2>&1 &
}

stop () {
  echo -n "dataloop agent: attempting deregister... "

  /usr/local/bin/dataloop-lin-agent --api-key <%= @api_key %> \
                                    --deregister

  if [ $? -ne 0 ]; then
    echo "failed! "
  else 
    echo "success"
  fi

  pkill -9 -f 'dataloop-lin-agent'
  echo "dataloop agent stopped"
}

case "$1" in
start)
        start
        ;;
stop)
        stop
        ;;
restart)
        stop
        sleep 1
        start
        ;;
*)
        echo "Usage: $0 { start | stop | restart }"
        exit 1
        ;;
esac

exit 0
